# system_script/crontab
0 8 * * * /bin/sh folder_creator.sh
@reboot python /path/to/aqua_botanical/program
@reboot bash /path/to/folder_creator/script